# supervisor.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:rmove]
command=rm /root/.ssh/*

[program:cmd]
command=/bin/bash

